# private_dot_config/systemd/user/atuin-server.service.d/override.conf
[Service]
# Atuin server configuration via environment variables
# All settings here override server.toml (which is not used)

# Network configuration
Environment=ATUIN_HOST=127.0.0.1
Environment=ATUIN_PORT=8888

# Database configuration
# Using SQLite with systemd %h specifier for dynamic home path
Environment=ATUIN_DB_URI=sqlite://%h/.local/share/atuin/server.db

# Registration (set to false after initial user setup)
Environment=ATUIN_OPEN_REGISTRATION=false

# Resource limits (increased for single-user self-hosted)
Environment=ATUIN_MAX_HISTORY_LENGTH=81920
Environment=ATUIN_MAX_RECORD_SIZE=2097152
Environment=ATUIN_PAGE_SIZE=2000
