# private_dot_config/systemd/user/atuin-server.service
[Unit]
Description=Atuin sync server (self-hosted)
Documentation=https://docs.atuin.sh/self-hosting/server-setup/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.nix-profile/bin/atuin server start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
ReadWritePaths=%h/.local/share/atuin %h/.config/atuin
CapabilityBoundingSet=
AmbientCapabilities=
NoNewPrivileges=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
PrivateTmp=true
PrivateDevices=true
LockPersonality=true

[Install]
WantedBy=default.target
